# scripts/letsencrypt-backup-tar.cron
0 15 * * * /nimbus/scripts/cert-update.sh > stdout
